# reverse-proxy/rootfs/etc/services.d/reverse-proxy/run
#!/usr/bin/with-contenv bash
exec /app/run.sh
