# valetudo-proxy/rootfs/etc/services.d/valetudo-proxy/run
#!/usr/bin/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash

# Run valetudo_proxy
URL=$(bashio::config 'url')
bashio::log.info "Valetudo proxy starting with target URL: ${URL}"
BASIC_AUTH=$(bashio::config 'basic_auth')
USERNAME=$(bashio::config 'username')
PASSWORD=$(bashio::config 'password')

export URL 
export BASIC_AUTH 
export USERNAME 
export PASSWORD

exec /usr/local/bin/valetudo-proxy